# systemd/daystill.timer
[Unit]
Description=Run DaystillRemedhan poster every 72 hours

[Timer]
# Start 5 minutes after boot and then every 72 hours (3 days)
OnBootSec=5min
OnUnitActiveSec=72h
# Persist across reboots so missed runs are triggered when system is back up
Persistent=true

[Install]
WantedBy=timers.target
